# linux/memduinod.service
[Unit]
Description=memduino daemon systemd service

[Service]
Type=simple
ExecStart=/usr/bin/memduinod

[Install]
WantedBy=multi-user.target